# deploy/guangzhou-github-proxy.service
[Unit]
Description=GitHub Proxy Service for VioletTeam
Documentation=https://github.com/your-repo
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/github-proxy
ExecStart=/usr/bin/python3 /opt/github-proxy/app.py
Restart=always
RestartSec=10

# 环境变量配置
# 代理设置（根据实际情况修改）
Environment="HTTP_PROXY=http://127.0.0.1:8118"
Environment="HTTPS_PROXY=http://127.0.0.1:8118"

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=github-proxy

# 安全配置
NoNewPrivileges=true
PrivateTmp=true

# 资源限制
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

